# systemd/sudocr.service
[Unit]
Description=SudOCR Application
Documentation=https://
